# config/systemd/user/simple_gpu_info.timer
[Unit]
Description=Run

[Timer]
OnActiveSec=60min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
